# luci-app-softethervpn/root/etc/uci-defaults/luci-app-softethervpn
#!/bin/sh
uci -q batch <<-EOF >/dev/null
	delete ucitrack.@softethervpn[-1]
	add ucitrack softethervpn
	set ucitrack.@softethervpn[-1].init=softethervpn
	commit ucitrack
EOF
exit 0
